# vps-worker/flowtix-worker.service
[Unit]
Description=Flowtix Facebook Automation Worker
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=__WORKER_DIR__
EnvironmentFile=-__WORKER_DIR__/.env
ExecStart=/usr/bin/node __WORKER_DIR__/worker.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/flowtix-worker.log
StandardError=append:/var/log/flowtix-worker.log

[Install]
WantedBy=multi-user.target
